# examples/systemd/voice-hotkey.service
[Unit]
Description=Voice hotkey daemon (Whisper + Hyprland)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.venvs/voice/bin/python <REPO_DIR>/voice-hotkey.py --daemon
Restart=on-failure
RestartSec=1
Environment=VOICE_AUDIO_BACKEND=pulse
Environment=VOICE_AUDIO_SOURCE=default
Environment=VOICE_COMMAND_MODEL=large-v3
Environment=VOICE_COMMAND_MODEL_EN=large-v3
Environment=VOICE_DICTATION_INJECTOR=wtype
Environment=VOICE_ASR_BACKEND=faster_whisper
# Environment=VOICE_ASR_BACKEND=whispercpp_server
# Environment=VOICE_WHISPER_SERVER_URL=http://127.0.0.1:8080/inference
Environment=VOICE_WAKE_SESSION_MAX_SECONDS=8
Environment=VOICE_WAKE_START_SPEECH_TIMEOUT_MS=7000
Environment=VOICE_WAKE_VAD_RMS_THRESHOLD=120
Environment=VOICE_WAKE_VAD_MIN_SPEECH_MS=20
Environment=VOICE_WAKE_VAD_END_SILENCE_MS=300

[Install]
WantedBy=default.target
